# packages/opencoat-runtime/systemd/opencoat-runtime.service
[Unit]
Description=OpenCOAT Runtime daemon
Documentation=https://github.com/HyperdustLabs/OpenCOAT
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/opencoat-daemon --config /etc/opencoat/config.yaml
Restart=on-failure
RestartSec=5
User=opencoat
Group=opencoat

# Hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/opencoat /run

[Install]
WantedBy=multi-user.target
